# contrib/systemd/arch-audit.timer
[Unit]
Description=Daily audit of installed packages against known vulnerabilities

[Timer]
OnCalendar=daily
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
